# recipes-custom/custom-domoticz-scripts/files/edftempo.cron
MAILTO=domoticz
10 06 * * * domoticz /etc/domoticz/edftempo.sh 2>&1 >> /dev/null
