# recipes-custom/custom-domoticz-scripts/files/edfwifinfo.cron
MAILTO=domoticz
* * * * * domoticz  /etc/domoticz/edfwifinfo.sh 2>&1 >> /dev/null
